# odysseus-ui.service
# Copy to /etc/systemd/system/odysseus-ui.service (install-service.sh does this for you)
[Unit]
Description=Odysseus UI
After=network.target

[Service]
Type=simple
User=david
WorkingDirectory=/home/david/odysseus
# PATH must include /home/david/.local/bin so the Browser MCP server can find `npx`
# (Node was installed there). Without it, the built-in Browser MCP fails to start.
Environment=PATH=/home/david/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/home/david/odysseus/venv/bin/uvicorn app:app --port 7000 --host 0.0.0.0
Restart=always
RestartSec=3
EnvironmentFile=-/home/david/odysseus/.env

[Install]
WantedBy=multi-user.target
